# macos/app.hexecute.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>app.hexecute</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/Hexecute.app/Contents/MacOS/hexecute</string>
    </array>
    <!-- Start at login and keep the resident agent alive so the global hotkey
         is always available. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <!-- Runs in the user's GUI session so it can create the overlay window. -->
    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
